# dehydrated/dehydrated.cron
#!/usr/bin/env bash

if [ -x /usr/sbin/%PRGNAM% ]; then
    /usr/sbin/%PRGNAM% --cron | logger -s -t %PRGNAM%
fi

